# src/thermoclock.service
# vim:ft=config:
[Unit]
Description=Thermoclock

[Service]
ExecStart=/usr/bin/thermoclock

[Install]
WantedBy=default.target

